# legacy/roles/roofcam/templates/roofcam-snapshot.cron
*/10 * * * * {{roofcam_dir}}/snapshot.sh